# modules/dev/systemd/dotfiles-watch.service
[Unit]
Description=Dotfiles watcher — auto-commit via inotify
After=network.target

[Service]
Type=simple
ExecStart=%h/dev/dotfiles/scripts/dotfiles-watch.sh
Restart=always
RestartSec=5
Nice=19
IOSchedulingClass=idle

[Install]
WantedBy=default.target
