# home/.config/systemd/user/boveda-auto-sync.service
[Unit]
Description=Auto-commit and push Boveda changes

[Service]
Type=oneshot
ExecStart=/var/home/alejandro/dev/dotfiles/scripts/boveda-auto-sync.sh
