# home/.config/systemd/user/dotfiles-auto-sync.service
[Unit]
Description=Auto-commit and push dotfiles changes

[Service]
Type=oneshot
ExecStart=/var/home/alejandro/dev/dotfiles/scripts/dotfiles-auto-sync.sh
