# home/.config/systemd/user/dotfiles-auto-sync.path
[Unit]
Description=Watch dotfiles for changes and auto-sync

[Path]
PathModified=/var/home/alejandro/dev/dotfiles
Unit=dotfiles-auto-sync.service

[Install]
WantedBy=default.target
